# stow/nvim-theme-mac/Library/LaunchAgents/com.peterszarvas.theme-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.peterszarvas.theme-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>export PATH="/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin"; bash "$HOME/.local/bin/mac-sync-nvim-theme" --apply; exec dark-notify -c "bash \"$HOME/.local/bin/mac-sync-nvim-theme\""</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/com.peterszarvas.theme-sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/com.peterszarvas.theme-sync.err.log</string>
</dict>
</plist>
