# stow/launchd/Library/LaunchAgents/com.peterszarvas.theme-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.peterszarvas.theme-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>export PATH="/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin"; "$HOME/.local/bin/mac-sync-theme" --apply</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>2</integer>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/com.peterszarvas.theme-sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/com.peterszarvas.theme-sync.err.log</string>
</dict>
</plist>
